# roles/system/templates/backup@.timer
[Unit]
Description=Backup timer %i

[Timer]
Unit=backup@%i.service
Persistent=false
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=default.target
