# roles/system/templates/backup@.service
[Unit]
Description=Backup %i

[Service]
Type=oneshot
ExecStart={{ deploy_path }}/%i/backup.sh
